# packages/node-exporter/files/DEBIAN/postinst
#!/bin/bash
systemctl enable node-exporter
systemctl start node-exporter
