# packages/prometheus/files/DEBIAN/postinst
#!/bin/bash
systemctl enable prometheus
systemctl start prometheus
